# ansible/40_thinkube/core/infrastructure/k8s/files/thinkube-node-ip-sync.service
# Managed by Thinkube — do not edit.
# Boot-time self-healing for the network-dependent node values; see
# /usr/local/bin/thinkube-node-ip-sync.
[Unit]
Description=Thinkube: sync kubelet node-ip to the current network
# Run after the network settles but always before kubelet registers.
Wants=network-online.target
After=network-online.target
Before=kubelet.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/thinkube-node-ip-sync
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
