# ansible/00_initial_setup/templates/vfio-bind.service.j2
# Copyright 2025 Alejandro Martínez Corriá and the Thinkube contributors
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=VFIO bind for PCI device {{ pci_slot }}
Documentation=https://github.com/thinkube/thinkube
DefaultDependencies=no
After=local-fs.target
After=systemd-modules-load.service
Before=sound.target
Before=display-manager.service
Before=graphical.target
# Run before the sound subsystem to claim audio devices first

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vfio-bind-{{ pci_slot | replace(':', '-') }}.sh
RemainAfterExit=yes
Restart=no
# Higher priority to override other kernel module loading
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=99
# Add root capabilities explicitly
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SYS_MODULE

[Install]
WantedBy=multi-user.target