# ansible/40_thinkube/core/infrastructure/k8s/templates/k8s-api-proxy.service.j2
# Node-local Kubernetes API VIP — proxy half (socket-activated).
# Forwards {{ k8s_stable_ip }}:6443 to the control plane's current LAN IP.
# This is the only place on a worker that references a network-dependent
# address; when the cluster moves networks, re-template this file and
# restart the service. Installed by 20_join_workers.yaml.
[Unit]
Description=Node-local Kubernetes API VIP proxy to control plane
Requires=k8s-api-proxy.socket
After=network-online.target

[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd --connections-max=512 {{ hostvars[control_plane_host]['lan_ip'] }}:6443
Restart=always
RestartSec=2
