# ansible/40_thinkube/core/infrastructure/k8s/templates/k8s-api-proxy.socket.j2
# Node-local Kubernetes API VIP — socket half.
# Listens on the same stable IP the control plane serves the API on, so
# every kubeconfig/agent on this worker can use the network-independent
# endpoint. Installed by 20_join_workers.yaml.
[Unit]
Description=Node-local Kubernetes API VIP socket ({{ k8s_stable_ip }}:6443)

[Socket]
ListenStream={{ k8s_stable_ip }}:6443
# Bind even if the k8s0 dummy interface is not up yet (no unit ordering
# against systemd-networkd needed).
FreeBind=yes

[Install]
WantedBy=sockets.target
